# private_dot_config/systemd/user/app-gaming.slice
[Slice]
MemorySwapMax=0M

[Service]
OOMScoreAdjust=1000
