# private_dot_config/systemd/user/fcitx5.service
[Unit]
Description=Fcitx5 for Niri

[Service]
Type=simple
ExecStart=fcitx5 -r

[Install]
WantedBy=niri.service
