# private_dot_config/systemd/user/set-day-theme.timer
[Unit]
Description=Set day theme at 6:00

[Timer]
OnCalendar=*-*-* 6:00:00
Persistent=true

[Install]
WantedBy=timers.target graphical.target
